# aegea/missions/arvados/rootfs.skel.in/etc/init/arvados-keepproxy.conf
description "Run Arvados Keep Proxy"

start on runlevel [2345]

respawn

exec keepproxy
