# aegea/missions/arvados/rootfs.skel.in/etc/init/arvados-keep.conf
description "Run Arvados Keep server"

start on runlevel [2345]

respawn

env GOGC=10
env GOMAXPROCS=4

exec keepstore -enforce-permissions=true -blob-signing-key-file=/etc/keepstore/blob-signing.key -volume=/mnt/keep
#-volume=/mnt2/keep
