# aegea/missions/arvados/rootfs.skel.in/etc/init/arvados-keepproxy.conf
description "Run Arvados Keep Proxy"

start on runlevel [2345]

respawn

env ARVADOS_API_TOKEN=$ARVADOS_API_TOKEN
env ARVADOS_API_HOST=$ARVADOS_UUID_PREFIX.$ARVADOS_HOSTNAME

exec keepproxy
