# aegea/missions/arvados/rootfs.skel.in/etc/init/arvados-nodemanager.conf
description "Run Arvados Node Manager"

start on runlevel [2345]

respawn

chdir /tmp

exec arvados-node-manager --config /etc/arvados/nodemanager/ec2.cfg
