# aegea/missions/arvados/rootfs.skel.in/etc/init/awslogs.conf
description "Run AWS CloudWatch Logs agent"

start on runlevel [2345]

respawn

env AWS_CONFIG_FILE=/etc/aws.conf

pre-start exec /bin/mkdir -p /var/awslogs/state

exec /usr/local/bin/aws logs push --config-file /etc/awslogs.conf
