# aegea/missions/arvados/rootfs.skel.in/etc/init/crunch-dispatch-slurm.conf
description "Run Arvados Crunch Dispatcher for SLURM"

start on runlevel [2345]

respawn

chdir /home/crunch
setuid crunch
setgid crunch

exec crunch-dispatch-slurm
