# aegea/rootfs.skel/etc/systemd/system/docker.service.d/ephemeralmount.conf
[Service]
ExecStartPre=/usr/bin/aegea-format-ephemeral-storage --yes
