# aegea/rootfs.skel/etc/systemd/system/docker.service.d/xb.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --log-driver=syslog --storage-driver=btrfs
