# aegea/rootfs.skel/etc/systemd/system/docker.service.d/xb.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --log-driver=awslogs --log-opt awslogs-group=docker --log-opt awslogs-region=us-west-2 --storage-driver=btrfs
