# aegea/missions/arvados/rootfs.skel.in/etc/profile.d/arvados.sh
export ARVADOS_API_HOST="$ARVADOS_UUID_PREFIX.$ARVADOS_HOSTNAME"
# export ARVADOS_API_TOKEN=
